# profile.d/zbar.sh
  export PATH='/app/vendor/node/bin:/app/bin:/app/vendor/bin:/app/node_modules/.bin:/app/vendor/bin'
  export LD_LIBRARY_PATH='/app/vendor/lib'